# crontab
# Domain Crawler Cron Jobs
# Format: minute hour day month day-of-week command
# Environment variables
PATH=/usr/local/bin:/usr/bin:/bin
PYTHONUNBUFFERED=1
PLAYWRIGHT_BROWSERS_PATH=/home/crawler/.cache/ms-playwright

# Run daily at 2 AM as crawler user with enhanced logging
0 2 * * * su -c 'cd /app && python domain_crawler.py --url-file /app/input/targets.txt --fqdn-list --concurrency 3 --output /app/output/domains.txt --log-file /app/logs/daily-$(date +\%Y\%m\%d).log --verbose' crawler

# Run weekly on Sunday at 3 AM for full validation with enhanced logging
#0 3 * * 0 su -c 'cd /app && python domain_crawler.py --url-file /app/input/targets.txt --fqdn-list --concurrency 5 --output /app/output/domains.txt --log-file /app/logs/weekly-$(date +\%Y\%m\%d).log --clear-cookies --verbose' crawler

# Example: Run every 6 hours
# 0 */6 * * * su -c 'cd /app && python domain_crawler.py --url-file /app/input/targets.txt --fqdn-list --concurrency 3 --output /app/output/domains.txt --log-file /app/logs/hourly-$(date +\%Y\%m\%d-\%H).log' crawler

# Note: Ensure this file ends with a newline
